# src/debian/postinst
#!/bin/sh

## Enables service
update-rc.d ardukey-auth defaults

exit 0
